# TSystems.LoveOTC/LoveOTC.service
[Unit]
Description=LoveOTC

[Service]
Type=notify
ExecStart=/home/aloento/LoveOTC/TSystems.LoveOTC

[Install]
WantedBy=multi-user.target
